# com.meltwater.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.meltwater.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/admin/Desktop/meltwater-automation/run_meltwater_sync.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/admin/Desktop/meltwater-automation</string>

    <!-- 测试阶段: 每天早上10点运行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- 正式阶段: 每周一早上10点运行 (启用时取消下面的注释,注释掉上面的 StartCalendarInterval) -->
    <!--
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    -->

    <key>StandardOutPath</key>
    <string>/Users/admin/Desktop/meltwater-automation/cron_output.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/admin/Desktop/meltwater-automation/cron_error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
